# config_files/supervisor/bot.conf
[program:bot]
directory = /opt/sms
command = /usr/local/bin/uwsgi -c /opt/sms/config.ini
autostart = true
autorestart = true
startsecs = 5
startretries = 30
user = root

