# config_files/supervisor/sms.conf
[program:sms]
directory = /opt/sms
command = /usr/bin/python /opt/sms/relay.py
autostart = true
autorestart = true
startsecs = 5
startretries = 30
user = root

